# claude-code-discord.service
# AI Bot - systemd service template
#
# Setup:
#   1. Copy this file:        sudo cp ai-bot.service /etc/systemd/system/
#   2. Edit the copy:         sudo systemctl edit --full ai-bot
#      - Set User, WorkingDirectory, and Environment paths for your system
#   3. Enable and start:      sudo systemctl enable --now ai-bot
#
# Logs:   sudo journalctl -u ai-bot -f
# Status: sudo systemctl status ai-bot

[Unit]
Description=AI Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# --- EDIT THESE FOR YOUR SYSTEM ---
User=YOUR_USERNAME
WorkingDirectory=/path/to/ai-bot
ExecStart=/bin/bash /path/to/ai-bot/run.sh
Environment=HOME=/home/YOUR_USERNAME
Environment=PATH=/home/YOUR_USERNAME/.deno/bin:/usr/local/bin:/usr/bin:/bin
# ----------------------------------
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
